# schedules\fixity-Fixity0.2-sch.xml
<?xml version="1.0" ?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
    <RegistrationInfo>
        <Date>2014-03-20T19:38:49</Date>
