# schedules/fixity-Fixity0.2-sch.xml
<?xml version="1.0" ?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
    <RegistrationInfo>
        <Date>2014-03-14T18:02:46</Date>
        <Author>Furqan</Author>
        <Version>Fixity Task Scheduler to Monitor A Folder Activity!</Version>
        <Description>Fixity Task Scheduler to Monitor A Folder Activity!</Description>
    </RegistrationInfo>
    <Triggers>
        <CalendarTrigger>
            <StartBoundary>2014-03-14T00:00:00</StartBoundary>
            <ScheduleByWeek>
                <WeeksInterval >1</WeeksInterval >
                <DaysOfWeek>
                    <Tuesday/>
                </DaysOfWeek >
            </ScheduleByWeek>
        </CalendarTrigger>
    </Triggers>
    <Settings>
        <Enabled>true</Enabled>
        <AllowStartOnDemand>true</AllowStartOnDemand>
        <AllowHardTerminate>true</AllowHardTerminate>
        <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
        <StartWhenAvailable>false</StartWhenAvailable>
        <WakeToRun>true</WakeToRun>
    </Settings>
    <Actions>
        <Exec>
            <Command>"D:\python\Fixity Project\schedules\fixity-Fixity0.2.vbs"</Command>
        </Exec>
      </Actions>
</Task>
