# assets/template/SchedulerWinMonthAdmin.xml
<?xml version="1.0" ?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
    <RegistrationInfo>
        <Date>{{current_data_time}}</Date>
        <Author>{{user}}</Author>
        <Version>{{version}}</Version>
        <Description>Fixity Task Scheduler to Monitor A Folder Activity!</Description>
    </RegistrationInfo>
    <Triggers>
        <CalendarTrigger>
            <StartBoundary>{{start_boundary}}</StartBoundary>
            <ScheduleByMonth>
                <DaysOfMonth>
                    <Day>{{days_of_month}}</Day>
                </DaysOfMonth >
                <Months>
                <January/>
                <February/>
                <March/>
                <April/>
                <May/>
                <June/>
                <July/>
                <August/>
                <September/>
                <October/>
                <November/>
                <December/>
                </Months>
            </ScheduleByMonth>
        </CalendarTrigger>
    </Triggers>
	<Principals> 
		<Principal id="Author">
			<UserId>SYSTEM</UserId>
			<RunLevel>HighestAvailable</RunLevel>
		</Principal>
	</Principals>
    <Settings>
        <Enabled>{{enabled}}</Enabled>
        <MultipleInstancesPolicy>Parallel</MultipleInstancesPolicy>
        <AllowStartOnDemand>{{allow_start_on_demand}}</AllowStartOnDemand>
        <AllowHardTerminate>{{allow_hard_terminate}}</AllowHardTerminate>
        <DisallowStartIfOnBatteries>{{disallow_start_if_on_batteries}}</DisallowStartIfOnBatteries>
        <StartWhenAvailable>{{start_when_available}}</StartWhenAvailable>
        <WakeToRun>{{wake_to_run}}</WakeToRun>
    </Settings>
    <Actions>
        <Exec>
            <Command>{{command}}</Command>
        </Exec>
      </Actions>
</Task>
